# voice-transcriber.service
[Unit]
Description=Voice Transcriber Service with Wayland/Hyprland Support
Documentation=Voice-to-text application using OpenAI Whisper
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/home/kdresdell/Documents/DEV/whisper/bin/voice-transcriber-service.sh
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitInterval=30

# Environment variables for Wayland/Hyprland
Environment=XDG_SESSION_TYPE=wayland
Environment=QT_QPA_PLATFORM=wayland
Environment=GDK_BACKEND=wayland
Environment=SDL_VIDEODRIVER=wayland
Environment=CLUTTER_BACKEND=wayland

# Working directory
WorkingDirectory=/home/kdresdell/Documents/DEV/whisper

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-transcriber

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/kdresdell/Documents/DEV/whisper/log /home/kdresdell/Documents/DEV/whisper/rec /home/kdresdell/Documents/DEV/whisper/metrics
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=default.target